# machine-state/launchd/user-LaunchAgents/com.dotfiles.machine-state-apply-next-day.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.dotfiles.machine-state-apply-next-day</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/ymy/.dotfiles/Scripts/machine-state-apply-next-day.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>50</integer>
    </dict>
  </dict>
</plist>
