# macOS/machine-state/launchd/user-LaunchAgents/ai.openclaw.cache-clean-2200.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/ymy</string>
		<key>TMPDIR</key>
		<string>/var/folders/z9/lk9dmnss5xb6x8xd6ztd72nh0000gn/T/</string>
	</dict>
	<key>Label</key>
	<string>ai.openclaw.cache-clean-2200</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>/Users/ymy/.openclaw/cache-clean-2200.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/ymy/.cache/openclaw/launchd-cache-clean.err</string>
	<key>StandardOutPath</key>
	<string>/Users/ymy/.cache/openclaw/launchd-cache-clean.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>22</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
</dict>
</plist>
