# macOS/machine-state/launchd/user-LaunchAgents/ai.openclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>
    
    <key>Comment</key>
    <string>OpenClaw Gateway (v2026.5.7)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>1</integer>
    <key>Umask</key>
    <integer>63</integer>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/ymy/.openclaw/service-env/ai.openclaw.gateway-env-wrapper.sh</string>
      <string>/Users/ymy/.openclaw/service-env/ai.openclaw.gateway.env</string>
      <string>/opt/homebrew/opt/node/bin/node</string>
      <string>/opt/homebrew/lib/node_modules/openclaw/dist/index.js</string>
      <string>gateway</string>
      <string>--port</string>
      <string>18789</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/ymy/.openclaw</string>
    <key>StandardOutPath</key>
    <string>/Users/ymy/.openclaw/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ymy/.openclaw/logs/gateway.err.log</string>
  </dict>
</plist>
